# etc/bin/alignak-app
#!/usr/bin/env bash
# -*- coding: utf-8 -*-

# Copyright (c) 2015-2016:
#   Matthieu Estrada, ttamalfor@gmail.com
#
# This file is part of (AlignakApp).
#
# (AlignakApp) is free software: you can redistribute it and/or modify
# it under the terms of the GNU Affero General Public License as published by
# the Free Software Foundation, either version 3 of the License, or
# (at your option) any later version.
#
# (AlignakApp) is distributed in the hope that it will be useful,
# but WITHOUT ANY WARRANTY; without even the implied warranty of
# MERCHANTABILITY or FITNESS FOR A PARTICULAR PURPOSE.  See the
# GNU Affero General Public License for more details.
#
# You should have received a copy of the GNU Affero General Public License
# along with (AlignakApp).  If not, see <http://www.gnu.org/licenses/>.

### BEGIN INIT INFO
# Provides:          alignak-app
# Required-Start:    $remote_fs $syslog
# Required-Stop:     $remote_fs $syslog
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: alignak application notifier
# Description:       alignak-app is a notifier for Alignak suite.
### END INIT INFO

LIB_NAME=alignak_app
DAEMON=alignak-app

ROOT_DIR="$HOME/.local/$LIB_NAME"
BIN_DIR="$ROOT_DIR/bin"
USER_BIN=$HOME/bin

# Verify import of alignak_app and PyQt
PYBIN=""
PYQT=0
PY=0

python -c "import alignak_app" >/dev/null  2>&1

if [ $? -eq 0 ]; then
    PY=2
fi

python -c "__import__('PyQt4')" >/dev/null  2>&1

if [ $? -eq 0 ]; then
    PYQT=2
fi

if [[ "$PY" -eq 2 && "$PYQT" -eq 2 ]]; then
    PYBIN=python
else
    python3 -c "import alignak_app"
    if [ $? -eq 0 ]; then
        PY=3
    fi
    python3 -c "__import__('PyQt5')"
    if [ $? -eq 0 ]; then
        PYQT=3
    fi
    if [[ "$PY" -eq 3 && "$PYQT" -eq 3 ]]; then
        PYBIN=python3
    else
        echo "It seems that your version of Python and PyQt are not compatible or not installed !"
        exit 1
    fi
fi



if [ ! -L "$USER_BIN/$DAEMON" ]; then
    echo "Alignak-app verify installation..."
    # Verify installation
    echo "The command is not yet created..."
    echo "Alignak-App check folders..."
    if [ ! -d "$USER_BIN" ]; then
        echo "Create bin folder for application..."
        mkdir "$USER_BIN"
    else
        echo "Folder already exists."
    fi
    echo "Create [alignak-app] command..."
    ln -s "$BIN_DIR/$(basename "$0")" "$USER_BIN/$DAEMON"
    echo "Installation is ok !"
fi

usage() {
    cat << END

Usage: $(basename "$0") {start|stop|status|restart}

END
}

version=`"$PYBIN" -c "from alignak_app import __version__; print(__version__)"`
doc_url=`"$PYBIN" -c "from alignak_app import __doc_url__; print(__doc_url__)"`

do_start() {
    echo "Starting $DAEMON daemon... with $PYBIN"
    echo "Version $version"
    echo "HELP: visit $doc_url"
    "$PYBIN" "$BIN_DIR/$DAEMON.py" &
}

do_stop() {
    pid=`ps fu |grep "alignak-app.py"|grep -v "grep"|awk '{print $2}'`
    if [ ! -z "$pid" ]; then
        echo "PID: $pid"
        kill "$pid"
        echo "Stop $DAEMON daemon..."
        echo "OK"
    else
        echo "PID: $pid"
        echo "$DAEMON is not running !"
    fi
}

do_status() {
    pid=`ps fu |grep "alignak-app.py"|grep -v "grep"|awk '{print $2}'`
    if [ ! -z "$pid" ]; then
        echo "$DAEMON is running... (pid $pid)"
    else
        echo "$DAEMON is not running !"
        echo "Run $DAEMON start"
    fi

}


CMD=$1

case "$CMD" in
    start)
        do_start
    ;;
    stop)
        do_stop
    ;;
    restart)
        do_stop
        do_start
    ;;
    *)
        usage
        exit 1
esac
exit 0